# recipes-netopeer2/netopeer2/netopeer2/netopeer2.conf
[program:netopeer2]
command=/usr/share/netopeer2/netopeer2_start.sh
autostart=true
autorestart=true
