# recipes-netopeer2/netopeer2/netopeer2/netopeer2.service
[Unit]
Description=Test

[Service]
ExecStart=/usr/share/netopeer2/netopeer2_start.sh

[Install]
WantedBy=multi-user.target